# 01_Introduction to Embedded Linux/session1/Backup Service/task2_backup/backup.timer
[Unit]
Descreption=Timer service for backup

[Timer]
oncalendar=*-*-* 02:00:00
Presistent=true

[Install]
Wantedby=multi-user.target

