# 01_Introduction to Embedded Linux/session1/Backup Service/task2_backup/backup.service
[Unit]
Description=backup service
Wants=backup.timer

[Service]
User=root
ExecStart=/home/Desktop/scripts/backup.bash
